# cron/cron.hourly/run_hourly
#!/bin/sh

/usr/local/bin/python3 /data/tostock/instock/job/realtime_only_job.py
#mkdir -p /data/logs
#DATE=`date +%Y-%m-%d:%H:%M:%S`
#echo $DATE >> /data/logs/hourly.log
